# deploy/demon-beauty.service
[Unit]
Description=Demon Beauty Web Builder
After=network.target

[Service]
Type=simple
WorkingDirectory=@ROOT_DIR@
ExecStart=@ROOT_DIR@/start.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
